# config/supervisord.conf
;process managment for gunicorn
[unix_http_server]
file=./pids/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
pidfile=./pids/supervisord.pid
logfile=./log/supervisord.log    ; supervisord log file
logfile_maxbytes=1MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

;TODO
; Use unicorn herder https://github.com/gds-operations/unicornherder to get reload via SIGURS2
; guide with unicorn herder https://weluse.de/blog/zero-downtime-deployments-with-unicorn-and-supervisord.html

[program:ratmapunicorn]
command=/usr/local/bin/bundle exec unicorn -c /ratmap/config/unicorn.rb 
autostart=true
autorestart=true
directory=/ratmap
user = root

[program:ratmapSidekiq]
command=/usr/local/bin/bundle exec sidekiq -C config/sidekiq.yml
autostart=true
autorestart=true
directory=/ratmap
user = root
